# wifi-script/wifi.service
[Unit]
Description

[Service]
Type=oneshot
ExecStart=/home/root/wifi.sh

[Install]
WantedBy=multi-user.target
